# docker/crontab
*/1 * * * * /scripts/power.sh > /proc/1/fd/1 2>/proc/1/fd/2
*/1 * * * * /scripts/inverter.sh > /proc/1/fd/1 2>/proc/1/fd/2
