# docker/crontab
*/1 * * * * /scripts/sun_calc.sh > /proc/1/fd/1 2>/proc/1/fd/2
*/1 * * * * /scripts/power_meter.sh > /proc/1/fd/1 2>/proc/1/fd/2
